# configs/gunicorn_mysite.socket
[Unit]
Description=gunicorn socket for mysite project

[Socket]
ListenStream=/run/gunicorn/mysite_socket

[Install]
WantedBy=sockets.target
